# crontab
# hello-cron
* * * * * echo "Hello world"
* * * * * echo "Hello world 2" >> /var/log/cron.log 2>&1